# alerts/com.berndgansbacher.pharmatrack-alerts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd job: run the EMA alerts watcher once a day at 08:00.

  Install:
    cp alerts/com.berndgansbacher.pharmatrack-alerts.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.berndgansbacher.pharmatrack-alerts.plist

  Run it once immediately to test:
    launchctl start com.berndgansbacher.pharmatrack-alerts

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.berndgansbacher.pharmatrack-alerts.plist

  The first ProgramArgument points at the python.org 3.11 framework python
  (where openpyxl is installed). If `which python3` shows a different path on
  your machine, update it to match — the system /usr/bin/python3 does NOT have
  openpyxl.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.berndgansbacher.pharmatrack-alerts</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Library/Frameworks/Python.framework/Versions/3.11/bin/python3</string>
    <string>/Users/berndgansbacher/Downloads/pharmatrack_-recent-drug-approvals/alerts/ema_alerts.py</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>8</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/berndgansbacher/Downloads/pharmatrack_-recent-drug-approvals/alerts/alerts.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/berndgansbacher/Downloads/pharmatrack_-recent-drug-approvals/alerts/alerts.log</string>
</dict>
</plist>
